# infra/0xeeTerm-shill.service
[Unit]
Description=0xeeTerm — $0xEE Shill-as-a-Service
After=network.target

[Service]
Type=oneshot
User=debian
WorkingDirectory=/home/debian/0xeeAI
ExecStart=/home/debian/0xeeAI/venv/bin/python3 /home/debian/0xeeAI/0xeeTerm shill
StandardOutput=append:/home/debian/0xeeAI/logs/shill.log
StandardError=append:/home/debian/0xeeAI/logs/shill.log

[Install]
WantedBy=multi-user.target
